# infra/systemd/carbuyer-postgres.service
[Unit]
Description=CarBuyer Postgres (Docker)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/markbohaychuk/repos/CarBuyerAssistant/infra
ExecStart=/usr/bin/docker compose up -d --wait postgres
ExecStop=/usr/bin/docker compose stop postgres

[Install]
WantedBy=multi-user.target
