# infra/systemd/carbuyer-discoverer.timer
[Unit]
Description=Run CarBuyer auction discoverer every 6 hours
After=carbuyer-postgres.service

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Persistent=true
Unit=carbuyer-discoverer.service

[Install]
WantedBy=timers.target
