# infra/systemd/carbuyer-distiller.service
[Unit]
Description=CarBuyer auction distiller (oneshot)
After=carbuyer-postgres.service
Requires=carbuyer-postgres.service

[Service]
Type=oneshot
User=mark
WorkingDirectory=/home/mark/repos/CarBuyerAssistant
EnvironmentFile=/home/mark/repos/CarBuyerAssistant/.env
ExecStart=/home/mark/repos/CarBuyerAssistant/.venv/bin/python -m carbuyer.apps.auction_distiller
StandardOutput=journal
StandardError=journal
