# infra/systemd/carbuyer-source-watchdog.timer
[Unit]
Description=Run CarBuyer stale-source watchdog hourly
After=carbuyer-postgres.service

[Timer]
OnCalendar=*:0/60
# 17min offset so the watchdog doesn't run at the same instant as the
# ingester's 6h tick (which it monitors); reduces noise from a transient
# overlap where ingest is mid-run but hasn't committed yet.
RandomizedDelaySec=300
Persistent=true
Unit=carbuyer-source-watchdog.service

[Install]
WantedBy=timers.target
